# scripts/babeld/files/DEBIAN/postinst
#!/bin/bash
echo "interface eth0" > /etc/babeld.conf
systemctl enable babeld
systemctl start babeld
