# scripts/yggdrasil-iptunnel/yggdrasil-iptunnel.service
[Service]
ExecStartPost=/usr/local/sbin/yggdrasil-setup
