# scripts/yggdrasil-iptunnel/yggdrasil-iptunnel.service
[Service]
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
ExecStartPost=/usr/local/sbin/yggdrasil-setup
