# scripts/mesh.service
[Unit]
Description=Mesh Service
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/mesh

[Install]
WantedBy=multi-user.target